# bin/etc/systemd/system/rusefi.service
#
# sudo systemctl enable CMsport_ECU
#

[Unit]
Description=CMsport_ECU systemd

[Service]
ExecStart=/rw/startup_root.sh

[Install]
WantedBy=multi-user.target
